# init/ts3.conf
start on (runlevel [2345])
stop on (runlevel [016])

pre-start script
/etc/init.d/ts3 start
end script

post-stop script
/etc/init.d/ts3 stop
end script
